# profile.d/go.sh
# https://github.com/golang/go/wiki/GOPATH
export GOPATH="/usr/share/gocode:$HOME/.local/gocode"
export PATH="$GOPATH/bin:$PATH"
export GOMAXPROCS=1
