# profile.d/python.sh
alias pip="pip --user"
alias pip2="pip2 --user"
alias pip3="pip3 --user"
