# network/if-up.d/openntpd
#!/bin/sh

CONFIG="/etc/openntpd/ntpd.conf"

if [ "${METHOD}"X = "loopback"X ] || [ "${METHOD}"X = "none"X ]
then
	exit 0
fi

invoke-rc.d openntpd force-reload || true
