# systemd/system/restorecond.service.d/override.conf
[Service]
CapabilityBoundingSet=
ProtectSystem=true      # restorecond handles files in /etc
ProtectHome=false       # restorecond handles files in /root
