# configs/shell-config/profile.d/wall.sh
# By default, refuse to receive wall(1) messages from unprivileged users
# Set `mesg y` in your config if this is the desired behaviour.
mesg n
